# auto-grayscale/auto-grayscale.service
[Unit]
Description=auto-grayscale: toggle bedtime mode at sunset/sunrise
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart="@DAEMON_PATH@" daemon
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
